# deb/lib/systemd/system/etrial-davserver@.service
[Unit]
Description=WebDAV server for user %i

[Service]
User=nobody
ExecStart=/usr/bin/etrial-davserver %i

[Install]
WantedBy=machines.target
