# config/gocryptfs.path
# This unit file is part of the etrial package.
#
# The etrial web app writes to /tmp/crypt.key when the user provides a
# decryption key on the /encrypted login page. This unit detects that, and
# starts gocryptfs.service, which will consume crypt.key to mount /crypt.
# 
# TODO: ensure that /tmp is volatile.

[Unit]
Description=encrypted storage passphrase

[Path]
PathChanged=/tmp/crypt.key
Unit=gocryptfs.service

[Install]
WantedBy=multi-user.target
